# deploy/merv-serve.service
[Unit]
Description=Mervin/Mervis model-switching server (serve.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ace
Group=ace
WorkingDirectory=/home/ace/Desktop/prjx/merv
Environment=HOME=/home/ace
Environment=PATH=/home/ace/.local/bin:/usr/local/bin:/usr/bin:/bin
# This box is a shared/contended VPS (CPU only). Pin the in-process backend and
# keep the thread count low -- llama.cpp busy-waits, so oversubscribing thrashes.
Environment=MERV_LLAMA_BACKEND=inproc
Environment=MERV_THREADS=2
ExecStart=/home/ace/.local/bin/uv run serve.py
Restart=always
RestartSec=5
# Model load (+ first-run HF download) can take a while; give it room.
TimeoutStartSec=600
StandardOutput=append:/home/ace/Desktop/prjx/merv/logs/service.log
StandardError=append:/home/ace/Desktop/prjx/merv/logs/service.log

[Install]
WantedBy=multi-user.target
